# debian/postinst
#!/bin/bash
#
# galliumos-braswell postinst
#

systemctl daemon-reload
systemctl unmask galliumos-braswell
systemctl enable galliumos-braswell
systemctl start galliumos-braswell

exit 0
